# etc/cron.d/ais
# m h dom mon dow user  command
01 00 * * *   root    /etc/init.d/ais restart

